# packaging/linux/st-server.service
[Unit]
Description=st low-latency game-streaming server
Documentation=https://github.com/pulstart/server
# graphical-session.target is started by the desktop environment once the
# session bus is up. Running under it means PipeWire, PulseAudio, and
# xdg-desktop-portal are all reachable by this unit natively.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/st-server
Restart=on-failure
RestartSec=2

# Set ST_CAPTURE=kms|pipewire|x11 by hand (via `systemctl --user edit
# st-server`) only if you need to force a specific backend for debugging.

[Install]
WantedBy=graphical-session.target
